# docker/8.2/etc/cron.d/nwwsoi-controller
* * * * * sail cd /var/www/html && ./artisan schedule:run >>/var/log/cron-nwwsoi-controller.log 2>&1
